# env/.config/systemd/user/swww-start.service
[Unit]
Description=Start swww daemon

[Service]
ExecStart=/usr/bin/swww-daemon
