# env/.config/systemd/user/swww-stop.service
[Unit]
Description=Stop swww

[Service]
ExecStart=/usr/bin/pkill swww-daemon
